# monitoring/prometheus-exporters/pihole/pihole_exporter.service
[Unit]
Description=PiHole Exporter

[Service]
User=pihole_exporter
ExecStart=/etc/pihole/pihole-exporter.sh

[Install]
WantedBy=multi-user.target